# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
logfile=/dev/fd/1
logfile_maxbytes=0

[program:npm]
priority=0
command=/bin/bash -c "cd ./frontend && npm run dev"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:flask]
priority=1
command=/bin/bash -c "FLASK_APP=run.py flask run --host 0.0.0.0"
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true